# openapi-chroma.service
[Unit]
Description=openapi-chroma MCP server
After=network.target

[Service]
Type=simple
User=mcp
WorkingDirectory=/home/mcp
ExecStart=/home/mcp/.bun/bin/bun run main.ts serve --transport http --port 3000
Restart=on-failure
RestartSec=5

Environment=CHROMA_DB_PATH=/var/lib/openapi-chroma
Environment=OLLAMA_URL=https://ai.home.itsnotcam.dev
Environment=OLLAMA_MODEL=snowflake-arctic-embed2
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
