# openapi-chroma.service
[Unit]
Description=openapi-chroma MCP server
After=network.target

[Service]
Type=simple
User=openapi-chroma
WorkingDirectory=/opt/openapi-chroma
ExecStart=/opt/openapi-chroma/.venv/bin/python main.py serve --transport http --port 3000
Restart=on-failure
RestartSec=5

# Keep ChromaDB data across restarts
Environment=CHROMA_DB_PATH=/var/lib/openapi-chroma
Environment=OLLAMA_URL=https://ai.home.itsnotcam.dev
Environment=OLLAMA_MODEL=mxbai-embed-large

[Install]
WantedBy=multi-user.target
